# trunk/scanbdservice.conf
# scanbd - scanner button daemon
# $Id$
#
#  scanbd - KMUX scanner button daemon
#
#  Copyright (C) 2008 - 2012  Wilhelm Meier (wilhelm.meier@fh-kl.de)
#
#  This program is free software; you can redistribute it and/or modify
#  it under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 2 of the License, or
#  (at your option) any later version.
#
#  This program is distributed in the hope that it will be useful,
#  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  GNU General Public License for more details.
#
#  You should have received a copy of the GNU General Public License
#  along with this program; if not, write to the Free Software
#  Foundation, Inc., 675 Mass Ave, Cambridge, MA 02139, USA.
#

description	"Scanner button daemon Service"

start on (filesystem and started dbus)
stop on stopping dbus

respawn

# let this env-var point to the place where the real sane configs are
# strip down the /etc/sane/dll.conf to the net-backend

env SANE_CONFIG_DIR=/usr/local/etc/scanbd

pre-start script
	  rm -f /var/run/scanbd.pid
end script

exec /usr/local/bin/scanbd -f -c /usr/local/etc/scanbd/scanbd.conf
